# software/os_stub/etc/systemd/system/take_picture.timer
[Unit]
Description=Takes picts

[Timer]
OnBootSec=30s
OnUnitActiveSec=600s

[Install]
WantedBy=timers.target

